# main.cron
30 07 * * * python /home/itde1/EFIM/main.py